# deploy/systemd/ttseed-dashboard.service
[Unit]
Description=ttseed dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/ttseed
ExecStart=/opt/ttseed/.venv/bin/ttseed-dashboard --config /opt/ttseed/config.yaml
Restart=on-failure
RestartSec=5
Environment=PYTHONUNBUFFERED=1
PrivateTmp=true
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
