# deploy/systemd/ttseed-ingest.service
[Unit]
Description=ttseed ingest
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/ttseed
ExecStart=/opt/ttseed/.venv/bin/python /opt/ttseed/src/main.py ingest --config /opt/ttseed/config.yaml
Environment=PYTHONUNBUFFERED=1
PrivateTmp=true
NoNewPrivileges=true
