# deploy/systemd/ttseed-ingest.timer
[Unit]
Description=ttseed ingest timer

[Timer]
OnCalendar=*:0/45
Persistent=true

[Install]
WantedBy=timers.target
